# infra/systemd/core-daemon.service
[Unit]
Description=CORE Background Worker Daemon
After=network.target postgresql.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
WorkingDirectory=/opt/dev/CORE
Environment="PATH=/opt/dev/CORE/.venv/bin:/home/lira/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="VIRTUAL_ENV=/opt/dev/CORE/.venv"
ExecStart=/opt/dev/CORE/.venv/bin/core-admin daemon start
Restart=on-failure
RestartSec=10
WatchdogSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=core-daemon

[Install]
WantedBy=default.target
